# energy-csv-studio/rootfs/etc/services.d/energy-csv-studio/run
#!/usr/bin/with-contenv bash
exec bash /run.sh
