# dmarc-monitor/rootfs/etc/services.d/dmarc-monitor/run
#!/usr/bin/with-contenv bashio

bash /run.sh
